# rootfs/etc/supervisor/conf.d/postfix.conf
[program:postfix]
command=/usr/sbin/postfix start-fg
startsecs=0
autorestart=false
